# systemd/discchatbot.service
[Unit]
Description=Discord Chatbot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/discchatbot
EnvironmentFile=/home/ubuntu/discchatbot/.env
ExecStart=/home/ubuntu/venv/bin/python bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
